# infra/systemd/chitin-pending-watch.timer
[Unit]
Description=Trigger chitin-pending-watch every 30s

[Timer]
OnBootSec=1min
OnUnitActiveSec=30s
Persistent=true

Unit=chitin-pending-watch.service

[Install]
WantedBy=timers.target
