# infra/systemd/chitin-shipped-entry-flipper.service
[Unit]
Description=Auto-flip backlog entries ready→partial when recent merged PRs shipped them
Documentation=file:%h/workspace/chitin/docs/runbooks/chitin-shipped-entry-flipper.md
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/workspace/chitin/scripts/shipped-entry-flipper.sh --apply

# Resource caps — tiny script
MemoryMax=256M
CPUQuota=100%
TimeoutStartSec=120

# Don't restart on failure — failure mode is git pull conflict or
# gh auth issue; auto-retry would just paper over it.
Restart=no

[Install]
WantedBy=default.target
