# infra/systemd/chitin-dispatcher.service
# Slice 7d: dispatcher one-shot, fired by chitin-dispatcher.timer.
#
# Reads docs/swarm-backlog.md, picks the next ready entry, submits one
# Temporal workflow, waits for completion, runs apply step. Idempotent:
# if a swarm/<entry-id>-* branch exists on origin OR a swarm-* workflow
# is currently running, this exits without dispatching.
#
# Install (paired with chitin-dispatcher.timer):
#   cp infra/systemd/chitin-dispatcher.{service,timer} ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now chitin-dispatcher.timer

[Unit]
Description=Chitin swarm dispatcher — one tick of backlog dispatch
# Need the worker up to actually process workflows we submit.
After=chitin-worker.service
Wants=chitin-worker.service

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/chitin
Environment=CHITIN_REPO_ROOT=%h/workspace/chitin
Environment=PATH=%h/.local/bin:%h/.vite-plus/bin:/snap/bin:/usr/local/bin:/usr/bin:/bin
# Optional per-user secrets (CHITIN_SLACK_WEBHOOK_URL etc). Leading dash
# = file is optional; absence is fine.
EnvironmentFile=-%h/.config/systemd/user/chitin.env
ExecStart=/home/red/.vite-plus/bin/pnpm exec tsx apps/temporal-worker/src/dispatcher.ts
StandardOutput=journal
StandardError=journal
# A single dispatch cycle should never run longer than this. If it does,
# something's wrong (hung workflow, stuck git, etc.) — kill the timer
# tick and let the next one retry.
#
# Must be > the largest tier wall_timeout in dispatcher.ts (currently
# T2/T3/T4 = 1800s) plus headroom for activity SIGKILL grace + apply
# step. 2400s = 40min gives ~10min headroom over 1800s wall.
TimeoutStartSec=2400
