# infra/systemd/chitin-codex-chain-ingest.service
# One-shot service fired by chitin-codex-chain-ingest.timer.
# Periodically pulls codex CLI's chain (~/.codex/sessions/) into
# chitin's analysis pipeline so codex tool calls show up alongside
# claude-code / copilot / gemini events in the unified chain.

[Unit]
Description=Chitin codex chain ingest

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/chitin
Environment=CHITIN_REPO_ROOT=%h/workspace/chitin
# `python3 -m analysis.codex_mine` requires the repo's `python/` dir
# on sys.path. WorkingDirectory is the repo root, but Python's
# default module search doesn't walk into `python/` — set
# PYTHONPATH explicitly. Without this the service fails with
# "ModuleNotFoundError: No module named 'analysis'".
Environment=PYTHONPATH=%h/workspace/chitin/python
Environment=PATH=%h/.local/bin:%h/.vite-plus/bin:/snap/bin:/usr/local/bin:/usr/bin:/bin
# Optional per-user secrets (CHITIN_SLACK_WEBHOOK_URL etc).
EnvironmentFile=-%h/.config/systemd/user/chitin.env
# Placeholder ExecStart — the matching `analysis.codex_mine` module
# is shipped by a paired backlog entry. Until it lands, this service
# will start and fail with "ModuleNotFoundError: No module named
# 'analysis.codex_mine'", which is intentional per the backlog
# entry's note (the unit is text pointing at a not-yet-existing
# path; tick exits non-zero, alarm-feeder eventually surfaces it).
ExecStart=/usr/bin/python3 -m analysis.codex_mine ingest --out-dir %h/.chitin
StandardOutput=journal
StandardError=journal
# 1h is enough for chain ingest on a typical day's codex traffic.
TimeoutStartSec=3600
