# infra/systemd/chitin-lessons.timer
# Timer for chitin-lessons.service.
# Daily — fires once a day, persistent so a missed tick from a powered-
# off rig fires on the next boot. The cadence is daily because lessons
# only meaningfully change on PR-merge events, and a daily extractor
# tick produces a tight digest the dispatcher can lean on. If the
# operator wants near-realtime lessons (rare), CHITIN_LESSONS_SCAN_LIMIT
# + a manual `systemctl --user start chitin-lessons.service` works.

[Unit]
Description=Chitin swarm lessons-learned timer

[Timer]
# Fire once per day. OnBootSec offsets so a fresh boot doesn't hammer
# the gh API alongside the daily-rollup tick.
OnUnitActiveSec=24h
OnBootSec=30min
Persistent=true
Unit=chitin-lessons.service

[Install]
WantedBy=timers.target
