# infra/systemd/chitin-shipped-entry-flipper.timer
[Unit]
Description=Trigger chitin-shipped-entry-flipper hourly
Documentation=file:%h/workspace/chitin/docs/runbooks/chitin-shipped-entry-flipper.md

[Timer]
# Hourly — recent merged PRs are the input; hourly cadence catches
# operator-merged PRs within an hour of merge.
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

Unit=chitin-shipped-entry-flipper.service

[Install]
WantedBy=timers.target
