# infra/systemd/chitin-pending-watch.service
[Unit]
Description=Poll hermes for operator approval replies; resolve pending_approvals

[Service]
Type=oneshot
ExecStart=%h/.local/bin/chitin-kernel pending watch-hermes
StandardOutput=journal
StandardError=journal

# Cap — single hermes shell-out + sqlite writes; no build, no network
# beyond hermes-gateway.
MemoryMax=128M
CPUQuota=50%
TimeoutStartSec=30

# Don't auto-restart — failures want operator visibility via the journal.
Restart=no
