# infra/systemd/chitin-swarm-rollup.timer
# Daily swarm health rollup — fires at 06:00 local time.
# Install alongside chitin-swarm-rollup.service:
#   cp infra/systemd/chitin-swarm-rollup.{service,timer} ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now chitin-swarm-rollup.timer

[Unit]
Description=Chitin swarm daily health rollup

[Timer]
# Fire at 06:00 local time every day.
OnCalendar=*-*-* 06:00:00
# If the system was off at 06:00, run once on next boot.
Persistent=true
Unit=chitin-swarm-rollup.service

[Install]
WantedBy=timers.target
