# infra/systemd/chitin-alarm-feeder.service
# One-shot service fired by chitin-alarm-feeder.timer.
# Reads the latest swarm-rollup JSON, dedups alarms against existing
# investigate-* backlog entries, drafts up to N (default 1) in_design
# entries with role:researcher per new alarm. Closes §7's "alarm
# fires → researcher investigates" loop.

[Unit]
Description=Chitin swarm alarm-feeder (rollup → backlog flywheel)
# Order after the rollup so it reads the freshest JSON.
After=chitin-swarm-rollup.service

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/chitin
Environment=CHITIN_REPO_ROOT=%h/workspace/chitin
Environment=PATH=%h/.local/bin:%h/.vite-plus/bin:/snap/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-%h/.config/systemd/user/chitin.env
ExecStart=/home/red/.vite-plus/bin/pnpm exec tsx apps/runner/src/alarm-feeder.ts
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60
