# infra/systemd/chitin-codex-usage-feed.service
# One-shot service fired by chitin-codex-usage-feed.timer.
# Re-mines ~/.codex/sessions/**/*.jsonl into the universal usage
# feed at ~/.cache/chitin/usage/codex.json. chitin-budget reads
# this for the Quota/usage section; chitin-status will read it
# for the per-driver dashboard row.
#
# Cheap: reads ~1MB of JSONL, writes ~1KB. Safe to run frequently.

[Unit]
Description=Refresh codex usage feed (rate_limits + call counts)

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/chitin
Environment=PATH=%h/.local/bin:%h/.vite-plus/bin:/snap/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-%h/.config/systemd/user/chitin.env
ExecStart=/bin/sh -c '\
  mkdir -p %h/.cache/chitin/usage && \
  cd %h/workspace/chitin/python/analysis && \
  uv run --project . python -m analysis.codex_mine usage \
    --write-feed=%h/.cache/chitin/usage/codex.json'
StandardOutput=journal
StandardError=journal
TimeoutStartSec=30
