# infra/systemd/chitin-researcher.timer
# Systemd timer for chitin-researcher.service
# Fires every 4 hours, persistent across reboots.

[Unit]
Description=Chitin swarm periodic researcher timer

[Timer]
OnUnitActiveSec=4h
OnBootSec=15min
Persistent=true
Unit=chitin-researcher.service

[Install]
WantedBy=timers.target
