# infra/systemd/chitin-pr-event-ingester.timer
# Generated by tools/generators/systemd-unit/generate.sh
# Fires chitin-pr-event-ingester.service on the configured schedule.

[Unit]
Description=Chitin PR event ingester (scans open PRs, dispatches review-graph) timer

[Timer]
OnUnitActiveSec=5min
OnBootSec=5min
# If Persistent=true, a missed tick fires on next boot.
Persistent=false
Unit=chitin-pr-event-ingester.service

[Install]
WantedBy=timers.target
