# infra/systemd/chitin-watchdog.service
[Unit]
Description=Chitin signal watchdog (locked agents, failed units, rollup alarms)

[Service]
Type=oneshot
ExecStart=%h/workspace/chitin/scripts/chitin-watchdog.sh

# notify-send needs DBus session access. systemd --user inherits
# DBUS_SESSION_BUS_ADDRESS from the user manager environment, so this
# works post-login. If notifications don't appear after install,
# inspect: systemctl --user show-environment | grep DBUS

# Caps — trivially small workload.
MemoryMax=64M
CPUQuota=30%
TimeoutStartSec=30

# Don't auto-restart — failures want operator attention via the log.
Restart=no

[Install]
WantedBy=default.target
