# infra/systemd/chitin-groomer.timer
# Timer for chitin-groomer.service.
# Daily — drips at most 1 candidate from arxiv into the backlog per
# day (cap configurable via CHITIN_GROOMER_CAP env). Persistent so a
# missed tick from a powered-off rig fires on the next boot. 1h boot
# offset stages it after rollup / lessons / debt-curator.

[Unit]
Description=Chitin swarm groomer timer

[Timer]
OnUnitActiveSec=24h
OnBootSec=1h
Persistent=true
Unit=chitin-groomer.service

[Install]
WantedBy=timers.target
