# infra/systemd/clawta-judge-backfill.timer
[Unit]
Description=Run Clawta swarm judge backfill every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true
Unit=clawta-judge-backfill.service

[Install]
WantedBy=timers.target
