# swarm/systemd/architecture-audit.timer
[Unit]
Description=Weekly architecture audit — Sunday 6am ET
Documentation=https://github.com/chitinhq/chitin/blob/main/swarm/bin/architecture-audit

[Timer]
# Fire weekly on Sunday at 6am America/Detroit. Persistent=true means a
# missed run (laptop asleep, restart, etc.) catches up on next available
# wakeup rather than being silently skipped.
OnCalendar=Sun *-*-* 06:00:00 America/Detroit
Persistent=true
AccuracySec=5min
Unit=architecture-audit.service

[Install]
WantedBy=timers.target
