# swarm/systemd/architecture-audit.service
[Unit]
Description=Weekly architecture audit — claude --print run
Documentation=https://github.com/chitinhq/chitin/blob/main/swarm/bin/architecture-audit
After=network-online.target hermes-gateway.service openclaw-gateway.service

[Service]
Type=oneshot
# The script needs:
#   - read access to ~/workspace/* (chitin + sibling repos)
#   - write access to a fresh audit worktree under ~/.cache/chitin/audit-worktrees/
#   - write access to ~/.openclaw/logs (script log)
#   - kanban CLI access (~/.hermes/kanban) for ticket creation
#   - gh CLI auth (~/.config/gh)
#   - claude config + cache for the headless run
ExecStart=%h/.local/bin/architecture-audit
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-%h/.config/chitin/architecture-audit.env
TimeoutStartSec=1800
Nice=15

# Hardening — read-mostly, with explicit write paths.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.openclaw/logs %h/.cache/chitin %h/workspace/chitin %h/.hermes/kanban %h/.config/gh %h/.claude %h/.cache/claude
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
