# swarm/systemd/clawta-poller.timer
[Unit]
Description=Clawta autonomous kanban dispatch poller — every 2 minutes
Documentation=https://github.com/chitinhq/chitin/blob/main/swarm/bin/clawta-poller

[Timer]
# Fire 30s after boot, then every 2min thereafter. RandomizedDelaySec
# spreads load if multiple swarm boxes ever share a board.
OnBootSec=30
OnUnitActiveSec=2min
RandomizedDelaySec=15s
AccuracySec=10s
Unit=clawta-poller.service
Persistent=true

[Install]
WantedBy=timers.target
