# swarm/systemd/clawta-poller.service
[Unit]
Description=Clawta autonomous kanban dispatch poller (one-shot tick)
Documentation=https://github.com/chitinhq/chitin/blob/main/swarm/bin/clawta-poller
After=network-online.target hermes-gateway.service openclaw-gateway.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/clawta-poller --once --max-dispatch 2
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
TimeoutStartSec=300
Nice=10

# Hardening — script only reads kanban DB + spawns clawta CLI; no need
# for network or write access outside ~/.openclaw/logs.
ProtectSystem=strict
ReadWritePaths=%h/.openclaw %h/.hermes/kanban
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true

# Make this user-service-only; never run as root.
[Install]
WantedBy=default.target
