# swarm/systemd/swarm-audit.service
[Unit]
Description=Daily swarm audit — operator-controlled (claude-code) summary
Documentation=https://github.com/chitinhq/chitin/blob/main/swarm/bin/swarm-audit
After=network-online.target hermes-gateway.service openclaw-gateway.service

[Service]
Type=oneshot
# The script needs:
#   - kanban sqlite read
#   - chain ledger read (~/.chitin/gov-decisions-*.jsonl)
#   - gh CLI for PR queries
#   - openclaw CLI for the agent call and the Discord post
# Set SWARM_AUDIT_DISCORD_CHANNEL=<id> in the user environment (e.g. via
# EnvironmentFile) to deliver to Discord; without it the audit logs the
# summary locally and exits — useful for first-run verification.
ExecStart=%h/.local/bin/swarm-audit
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-%h/.config/chitin/swarm-audit.env
TimeoutStartSec=300
Nice=10

# Hardening — reads several user-state directories, writes only to logs.
ProtectSystem=strict
ReadWritePaths=%h/.openclaw/logs
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
