# python/argus/systemd/argus-indexer.service
[Unit]
Description=Argus Observatory Indexer
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=%h
ExecStart=/usr/bin/env python -m argus index --decisions-dir %h/.chitin --follow
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
